# dist/deb/DEBIAN/postinst
#!/bin/bash

systemctl enable mako-server
systemctl start mako-server
