# dist/deb/DEBIAN/prerm
#!/bin/bash

systemctl disable mako-server
systemctl stop mako-server
